# dot_config/systemd/user/dv.service.tmpl
[Unit]
Description=dv(0x5341/dv)
After=network-online.target

[Service]
Environment=DV_SERVER_HOST=0.0.0.0
Environment=DV_SERVER_PORT=3000
Environment=DV_VIBE_KANBAN_PORT=4000
Environment=DV_GHQ_ROOT={{.chezmoi.homeDir}}/ghq
Environment=DV_GHQ_PATH={{.chezmoi.homeDir}}/.local/share/mise/shims/ghq
Environment=PATH={{.chezmoi.homeDir}}/.local/share/mise/shims:{{.chezmoi.homeDir}}/.local/share/aquaproj-aqua/bin
Environment=AQUA_GLOBAL_CONFIG={{.chezmoi.homeDir}}/.config/aquaproj-aqua/aqua.yaml
Environment=AQUA_POLICY_CONFIG={{.chezmoi.homeDir}}/.config/aquaproj-aqua/aqua-policy.yaml
ExecStart={{.chezmoi.homeDir}}/.local/share/aquaproj-aqua/bin/dv
Restart=always
RuntimeMaxSec=86400

[Install]
WantedBy=default.target
